# supervisord.conf
[supervisord]
nodaemon=true

[program:mongodb]
command=/run.sh  -DFOREGROUND #To adapt, I don't know how to launch your mongodb process

[program:apache2]
command=/usr/sbin/rsyslogd -D
